# python-transcription-service/transcribe/laravel.conf
[program:transcribe]
process_name=%(program_name)s_%(process_num)02d
command=php /srv/html/transcribe/artisan queue:work --queue=transcribe --max-time=900 --sleep=5 --tries=3 --timeout=850 --delay=300
autostart=true
autorestart=true
user=www-data
numprocs=8
startretries=1000000
redirect_stderr=true
stdout_logfile=/srv/html/transcribe/storage/logs/supervisord_queue-%(program_name)s_%(process_num)02d.log
stdout_logfile_maxbytes=1000000
stdout_logfile_backups=10
stopasgroup=true
killasgroup=true
